# src/agent-th.service.in
[Unit]
Description=Temperature & Humidity
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=nobody
Group=dialout
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
ExecStart=@prefix@/bin/agent-th
Restart=always

[Install]
WantedBy=bios.target
