# src/bios-agent-th.service.in
[Unit]
Description=Temperature & Humidity
After=malamute.service
Requires=malamute.service
PartOf=bios.target

[Service]
Type=simple
# Note: This service requires access to serial ports, etc.
User=nobody
Group=dialout
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
ExecStart=@prefix@/bin/bios-agent-th
Restart=always

[Install]
WantedBy=bios.target
